# deploy/systemd/winimi-backend-queue.service
[Unit]
Description=Winimi Bakery Laravel queue worker
Documentation=https://github.com/sajadkhavas/winimi-bakery-backend/blob/main/docs/PHASE_19_PRODUCTION_DEPLOYMENT.md
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=winimi
Group=www-data
WorkingDirectory=/var/www/winimi/backend/current
ExecStart=/usr/bin/php artisan queue:work database --sleep=3 --tries=3 --backoff=5 --timeout=60 --max-time=3600
Restart=always
RestartSec=5s
TimeoutStopSec=90s
KillSignal=SIGTERM
SyslogIdentifier=winimi-backend-queue
UMask=0027

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/winimi/backend/shared/storage /var/www/winimi/backend/current/bootstrap/cache
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
