# deploy/systemd/winimi-backend-backup.timer
[Unit]
Description=Run Winimi Bakery production backup daily

[Timer]
OnCalendar=*-*-* 02:30:00
RandomizedDelaySec=10m
Persistent=true
Unit=winimi-backend-backup.service

[Install]
WantedBy=timers.target
