# deploy/systemd/winimi-backend-scheduler.service
[Unit]
Description=Winimi Bakery Laravel scheduler tick
Documentation=https://github.com/sajadkhavas/winimi-bakery-backend/blob/main/docs/PHASE_19_PRODUCTION_DEPLOYMENT.md
After=network-online.target

[Service]
Type=oneshot
User=winimi
Group=www-data
WorkingDirectory=/var/www/winimi/backend/current/app
ExecStart=/usr/bin/php artisan schedule:run --no-interaction
SyslogIdentifier=winimi-backend-scheduler
UMask=0027

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/winimi/backend/shared/storage /var/www/winimi/backend/current/app/bootstrap/cache
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictRealtime=true
